# multiContainer-paas/services/flask/var/spool/cron/crontabs/root
59 23 * * * /usr/sbin/logrotate -f /etc/logrotate.d/uwsgi
59 23 * * * /var/www/spdpaas/regularly/remove_files.sh /var/www/html/CHUNZU/csvfile all
0 * * * * /usr/local/bin/python /var/www/spdpaas/regularly/check_redis_celery_queue_length.py >> /var/www/spdpaas/regularly/check_redis_celery_queue_length.log 2>&1