# celery-worker/simple_worker/var/spool/cron/crontabs/root
25 0 * * * /var/app/use_env_run_python.sh
