# multiContainer-paas/services/php/var/spool/cron/crontabs/root
* * * * * php /var/www/html/API/mes_modules/mes_timer/machine_on_off_hist.php
0 8 * * * php /var/www/html/API/mes_modules/mes_timer/device_reset.php
2 8 * * * php /var/www/html/API/mes_modules/mes_timer/machine_status_sum.php
1 * * * * php /var/www/html/API/mes_modules/mes_timer/device_basic.php
2 8 */2 * * rm /var/www/html/API/mes_modules/newfile.txt