# multiContainer-paas/services/celery/var/spool/cron/crontabs/root
59 23 * * * /usr/sbin/logrotate -f /etc/logrotate.d/celery
25 0 * * * /usr/local/bin/python /var/app/dailyCalculateAggregatelyApirecords.py > /var/app/dailyCalculateAggregatelyApirecords_output.txt 2>&1
