# multiContainer-paas/services/celery/var/spool/cron/crontabs/root
25 0 * * * /usr/local/bin/python /var/app/dailyCalculateAggregatelyApirecords.py > /var/app/dailyCalculateAggregatelyApirecords_output.txt 2>&1
* * * * * echo "test" >> /var/app/dailyCalculateAggregatelyApirecords_output.txt